# dev-python/pyafipws/files/pyafipws.cron
#!/bin/bash
# Descarga y procesa el padrón desde AFIP.
#
# Copyright (C) 2004 - 2020 NIS
# Martin Andres Gomez Gimenez <mggimenez@nis.com.ar>
# Distributed under the terms of the GNU General Public License v3
#

# Comentar la siguiente línea para actualizar diariamente.
exit 1


PYTHONPATH="SET_PYTHONPATH"
CACHEPATH="${PYTHONPATH}/cache"

cd "${CACHEPATH}"
rm -f padron.txt
python3 "${PYTHONPATH}/padron.py" --descargar &>/dev/null

# Cambia codificación a UTF-8
mv padron.txt padron_orig.txt
iconv -f ISO-8859-15 -t UTF8 padron_orig.txt > padron.txt
rm padron_orig.txt

python3 "${PYTHONPATH}/padron.py" --borrar --procesar &>/dev/null
