# dev-python/pyafipws/files/pyafipws-2.7.20191021.cron
#!/bin/bash
# Descarga y procesa el padrón desde AFIP.
#
# Copyright (C) 2012 - 2016 Ingenio Virtual
# Martin Andres Gomez Gimenez <mggimenez@ingeniovirtual.com.ar>
# Distributed under the terms of the GNU General Public License v3
#


PYTHONPATH="/var/lib/odoo/.cache/pyafipws"
CACHEPATH="/var/lib/odoo/.cache/pyafipws"

cd "${CACHEPATH}"
python2 "${PYTHONPATH}/padron.py" --descargar --procesar &>/dev/null
